# .files
/home/michael/bin/brightness
/etc/rc.conf
/etc/rc.d/emacsd
