# .files
/home/michael/bin/brightness
/home/michael/bin/gitconf
/home/michael/bin/volume
/home/michael/.xmonad/xmonad.hs
/home/michael/.bashrc
/home/michael/.bash_aliases
/home/michael/.bash_profile
/home/michael/.xinitrc
/home/michael/.fonts.conf
/home/michael/bin/ifexping
/home/michael/bin/touchpad.sh
/home/michael/bin/translate.py
/home/michael/bin/translate.sh
/home/michael/.inputrc
/home/michael/.emacs.d/init.el
/home/michael/.emacs.d/lisp/zenburn.el
/home/michael/.emacs.d/lisp/linum-relative.el
/home/michael/dev/.main/linum-relative.el
/home/michael/bin/autorun3gmodem.sh
/etc/rc.d/rc.local
/home/michael/bin/slocktogglevolume.sh
/etc/rc.d/rc.transmission
/home/michael/bin/rc.d
/etc/rc.d/rc.6
/etc/rc.d/rc.K
/etc/rc.d/rc.M
/home/michael/bin/nmvpn_auto2
/home/michael/.pentadactyl/plugins/emacs.penta
/home/michael/.pentadactylrc
/home/michael/bin/screenshot
/home/michael/bin/md5crack
/home/michael/bin/git2pkg
/home/michael/bin/vksearch
/home/michael/bin/transmission-dog
/home/torrents/bin/torrentdownloadcomplete.sh
/home/michael/etc/goodies/settings.json.work3
/home/michael/etc/goodies/transmission.settings.json
/etc/apt/preferences
/etc/apt/sources.list
/home/michael/.gdbinit
/home/michael/.zshrc
/home/michael/.ext_aliases.zsh
/etc/init.d/saged
