# .files
/home/michael/bin/brightness
/etc/rc.conf
/etc/rc.d/emacsd
/home/michael/bin/gitconf
/home/michael/bin/volume
/home/michael/.xmonad/xmonad.hs
/home/michael/.bashrc
/home/michael/.bash_aliases
/home/michael/.bash_profile
/home/michael/.bash_logout
/home/michael/.xinitrc
/home/michael/.fonts.conf
/home/michael/bin/ifexping
/home/michael/bin/touchpad.sh
/home/michael/bin/translate.py
/home/michael/bin/translate.sh
/home/michael/.inputrc
/home/michael/.emacs.d/init.el
